# meta-str/recipes-str/tc-str-manager/tcc805x-cluster/poky-telechips-systemd/tc-str-manager.service
[Unit]
Description=telechips str manager for suspend & resume
After=pulseaudio.service

[Service]
Type=forking
PIDFile=/run/str-man.pid
ExecStart=/usr/bin/TCSTRManager -d 1 -t main 
[Install]
WantedBy=multi-user.target

